# dokploy/supervisord.conf
[supervisord]
nodaemon=true
logfile=/dev/null
logfile_maxbytes=0

[program:postgres]
command=postgres -D /var/lib/postgresql/data -h 127.0.0.1
user=postgres
autostart=true
autorestart=true
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
redirect_stderr=true
priority=10

[program:backend]
command=uvicorn app.main:app --host 127.0.0.1 --port 8000 --log-level info
directory=/app
environment=POSTGRES_HOST="127.0.0.1",POSTGRES_PORT="5432",POSTGRES_DB="turboea",POSTGRES_USER="turboea"
autostart=true
autorestart=true
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
redirect_stderr=true
startsecs=5
priority=20

[program:nginx]
command=nginx -g "daemon off;"
autostart=true
autorestart=true
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
redirect_stderr=true
priority=30
